# lookoukwindow.service
[Unit]
Description=Lookoukwindow - NASA 太空直播和 Google 相册展示
After=network.target

[Service]
Type=simple
User=jenkins
WorkingDirectory=/var/www/loooutwindow
Environment="PATH=/var/www/loooutwindow/venv/bin"
ExecStart=/var/www/loooutwindow/venv/bin/python /var/www/loooutwindow/run.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
